# etc/systemd/system/zabbix-haproxy.service
[Unit]
Description=Monitor HAProxy

[Service]
Type=simple
ExecStart=/etc/zabbix/scripts/zhaproxy.py

[Install]
WantedBy=multi-user.target

